# websocket-start.conf
[program:websocket-start]
directory=/var/www/TestTask/
command=php artisan websockets:serve
autostart=true
autorestart=true
